# gatesopener/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare base_url
declare token
declare key

base_url=$(bashio::config 'base_url')
token=$(bashio::config 'token')
key=$(bashio::config 'key')

bashio::log.info "Starting gates opener server with ${base_url}"

## Run your program
/usr/bin/gates_opener --token "${token}" --base-url "${base_url}" --internal-key "${key}"